# vps-hermes/systemd/hako-creative-worker.service
[Unit]
Description=HAKO Creative persistent job worker
After=network-online.target postgresql.service hako-creative-api.service
Wants=network-online.target

[Service]
Type=simple
User=hako-creative
Group=hako-creative
WorkingDirectory=/srv/hako-creative/current
EnvironmentFile=/etc/hako-creative/runtime.env
# O browser do HyperFrames precisa de cache persistente fora da release, que é
# somente leitura para o serviço. ProtectHome=true bloqueia o $HOME real, e
# /srv/hako-creative/cache já está em ReadWritePaths. Sem estas duas linhas o
# job render_prototype falha ao baixar/abrir o browser.
Environment=HOME=/srv/hako-creative/cache
Environment=XDG_CACHE_HOME=/srv/hako-creative/cache/.cache
ExecStart=/srv/hako-creative/current/.venv/bin/python tools/runtime_worker.py
Restart=on-failure
RestartSec=3s
TimeoutStopSec=30s
KillSignal=SIGTERM
UMask=0027

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/srv/hako-creative/assets /srv/hako-creative/cache /srv/hako-creative/logs

MemoryHigh=1536M
MemoryMax=2G
CPUQuota=180%
TasksMax=256
LimitNOFILE=8192

[Install]
WantedBy=multi-user.target
